# monolight.etc/systemd/system/monolight.service
[Unit]
Description=Monolight RGB matrix

[Service]
Type=oneshot
ExecStart=/usr/bin/monolight-led

[Install]
WantedBy=multi-user.target
